# aiticket-standalone/src/APP/deploy_scripts/supervisor.conf
[program:ai-ticket]
; AI工单系统后端服务配置
; 安装: sudo cp supervisor.conf /etc/supervisor/conf.d/ai-ticket.conf
; 重载: sudo supervisorctl reread && sudo supervisorctl update

directory=/opt/ai-ticket/APP/backend
command=/opt/ai-ticket/APP/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 18000
user=ubuntu
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=10
stopsignal=TERM
stderr_logfile=/var/log/supervisor/ai-ticket.err.log
stdout_logfile=/var/log/supervisor/ai-ticket.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PYTHONPATH="/opt/ai-ticket/APP/backend",AITICKET_HOST="qcl",AITICKET_PEER_BRIDGE="http://127.0.0.1:13800"
